# configs/supervisor_nginx.conf
[program:nginx]
command=nginx -c /Monitor/configs/nginx.conf
autostart=true
autorestart=true
user=root
redirect_stderr=false

[supervisord]
logfile=/dev/null